# arch/install
#!/bin/bash

# `-e` Exit on error.
# `-u` Treat unset variables as an error.
set -eu

# Determine the dotfiles directory.
DOTFILES_DIRECTORY="$(git rev-parse --show-toplevel)"

# Install Laravel development environment.
if [ -x "${HOME}/.local/share/omarchy/bin/omarchy-install-dev-env" ]; then
  "${HOME}/.local/share/omarchy/bin/omarchy-install-dev-env" laravel
fi

# Optimize resolution for my monitors.
if [ -f "${HOME}/.config/hypr/monitors.conf" ]; then
  sed -i \
    -e 's/^env = GDK_SCALE,2/# &/' \
    -e 's/^monitor=,preferred,auto,auto/# &/' \
    -e 's/^# env = GDK_SCALE,1$/env = GDK_SCALE,1/' \
    -e 's/^# monitor=,preferred,auto,1$/monitor=,preferred,auto,1/' \
    "${HOME}/.config/hypr/monitors.conf"
fi

# Allow to toggle between keyboard layouts.
if [ -f "${HOME}/.config/hypr/input.conf" ]; then
  sed -i \
    -e 's/^.*kb_layout.*/  kb_layout = us,es/' \
    -e 's/^\(.*kb_options = compose:caps\).*/\1,grp:alts_toggle/' \
    "${HOME}/.config/hypr/input.conf"
fi

# Show the active keyboard layout in waybar.
WAYBAR_CONFIG="${HOME}/.config/waybar/config.jsonc"
if ! jq -e '."modules-right" | index("hyprland/language")' "${WAYBAR_CONFIG}" >/dev/null 2>&1; then
  jq '."modules-right" += ["hyprland/language"] | . + {"hyprland/language": {"format": "{short}"}}' \
    "${WAYBAR_CONFIG}" >"${WAYBAR_CONFIG}.tmp" && mv "${WAYBAR_CONFIG}.tmp" "${WAYBAR_CONFIG}"
  omarchy-restart-waybar
fi

# Disable IDLE locking until the `hyprlock` bug is fixed upstream.
LINE="exec-once = while ! pgrep -x hypridle; do sleep 0.2; done && omarchy-toggle-idle"
grep -qF "$LINE" "${HOME}/.config/hypr/autostart.conf" 2>/dev/null || echo "$LINE" >>"${HOME}/.config/hypr/autostart.conf"

# Disable hardware rendered cursors to preven visual artifacts.
LINE="env = WLR_NO_HARDWARE_CURSORS,1"
grep -qF "$LINE" "${HOME}/.config/hypr/envs.conf" 2>/dev/null || echo "$LINE" >>"${HOME}/.config/hypr/envs.conf"

# Install AUR Packages.
applications=(
  jetbrains-toolbox
  lmstudio-bin
  slack-desktop
)
for application in "${applications[@]}"; do
  yay -S --noconfirm --needed "${application}"
done
unset application
unset applications

# Install the .NET SDK for C# development.
sudo pacman -S --noconfirm --needed dotnet-sdk

# Install the Circleback CLI for terminal access to meeting notes.
if ! command -v circleback >/dev/null 2>&1; then
  npm install -g @circleback/cli
fi

# Install GitHub Copilot CLI for terminal-based agentic coding.
if ! command -v copilot >/dev/null 2>&1; then
  npm install -g @github/copilot
fi

# Install Web Applications
web_applications=(
  "Claude|https://claude.ai/new|https://cdn.jsdelivr.net/gh/homarr-labs/dashboard-icons/png/claude-ai-light.png"
  "Todoist|https://app.todoist.com/app/today|https://cdn.jsdelivr.net/gh/homarr-labs/dashboard-icons/png/todoist.png"
)
for web_application in "${web_applications[@]}"; do
  IFS='|' read -r name url icon <<<"${web_application}"
  omarchy-webapp-install "${name}" "${url}" "${icon}"
done
unset web_application
unset web_applications

# Enable SSH daemon for incoming connections.
sudo systemctl enable --now sshd

# Enable SSH agent via systemd for passphrase caching.
systemctl --user enable --now ssh-agent.socket

# Set SSH_AUTH_SOCK for the SSH agent in Hyprland.
LINE='env = SSH_AUTH_SOCK,$XDG_RUNTIME_DIR/ssh-agent.socket'
grep -qF "$LINE" "${HOME}/.config/hypr/envs.conf" 2>/dev/null || echo "$LINE" >>"${HOME}/.config/hypr/envs.conf"

# Auto-add SSH keys to the agent on first use.
mkdir -p "${HOME}/.ssh" && chmod 700 "${HOME}/.ssh"
touch "${HOME}/.ssh/config" && chmod 600 "${HOME}/.ssh/config"
if ! grep -q "AddKeysToAgent" "${HOME}/.ssh/config"; then
  printf '\nHost *\n  AddKeysToAgent yes\n' >>"${HOME}/.ssh/config"
fi

# Source custom bash aliases.
ALIASES_SOURCE="source \"${DOTFILES_DIRECTORY}/arch/bash/aliases\""
grep -qF "${ALIASES_SOURCE}" "${HOME}/.bashrc" 2>/dev/null || echo "${ALIASES_SOURCE}" >>"${HOME}/.bashrc"

# Install stow for managing dotfile symlinks.
sudo pacman -S --noconfirm --needed stow

# Install lazygit configuration file.
stow -d "${DOTFILES_DIRECTORY}/arch" -t "${HOME}" lazygit

# Stow NeoVim configuration (hidden files, etc.).
stow -d "${DOTFILES_DIRECTORY}/arch" -t "${HOME}" nvim


# Apply Bluetooth audio fixes for the onboard WiFi/BT card.
"${DOTFILES_DIRECTORY}/arch/bluetooth-audio-rtl8821ce"

# Apply gaming fixes for this NVIDIA + Hyprland workstation.
"${DOTFILES_DIRECTORY}/arch/gaming-nvidia-hyprland"

LAZYVIM_JSON="${HOME}/.config/nvim/lazyvim.json"
if [ -f "${LAZYVIM_JSON}" ]; then
  jq '.extras |= (. + [
    "lazyvim.plugins.extras.dap.core",
    "lazyvim.plugins.extras.lang.dotnet",
    "lazyvim.plugins.extras.lang.php",
    "lazyvim.plugins.extras.lang.typescript",
    "lazyvim.plugins.extras.lang.vue",
    "lazyvim.plugins.extras.formatting.prettier",
    "lazyvim.plugins.extras.linting.eslint"
  ] | unique)' "${LAZYVIM_JSON}" > "${LAZYVIM_JSON}.tmp" && mv "${LAZYVIM_JSON}.tmp" "${LAZYVIM_JSON}"
fi
